# config/etc/profile.d/unset.sh
#!/usr/bin/env bash

unset ADMIN_PASS
unset ADMIN_USER
unset CEJ_PASS
unset CEJ_USER
unset MAX_PASS
unset MAX_USER
unset REP_PASS
unset REP_USER
unset S3_ACCESS_KEY_ID
unset S3_BUCKET
unset S3_ENDPOINT
unset S3_REGION
unset S3_SECRET_ACCESS_KEY
